# gauss_bringup/systemd/gauss-ros.service

Description="bringup gauss-ros"

[Service]
Type=simple
ExecStart=/usr/local/bin/gauss-ros-start

[Install]
WantedBy=multi-user.target
